# deploy/systemd/atelier-kernel.service
[Unit]
Description=DjinnOS Shygazun Kernel Service
After=network.target

[Service]
Type=simple
User=djinn
Group=djinn
WorkingDirectory=/opt/djinnos/DjinnOS-Shyagzun
ExecStart=/opt/djinnos/.venv/bin/python -m uvicorn shygazun.kernel_service:app --host 127.0.0.1 --port 8000 --app-dir /opt/djinnos/DjinnOS-Shyagzun
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
